# PiSetup/PyRemote.service
[Unit]
Description=PyRemote
After=multi-user.target
[Service]
Type=idle
ExecStart=/home/Stefaan/PyRemote/run.sh
[Install]
WantedBy=multi-user.target